# root/etc/s6-overlay/s6-rc.d/docker-conf/up
#!/usr/bin/with-contenv bash
exec /scripts/docker-conf.sh