# root/etc/s6-overlay/s6-rc.d/adduser/up
#!/usr/bin/with-contenv bash
exec /scripts/adduser.sh
