# root/etc/s6-overlay/s6-rc.d/banner/up
#!/usr/bin/with-contenv bash
exec /scripts/banner.sh
